# distribution/macos/postinstall
#!/bin/sh

set -e
/bin/launchctl load "/Library/superplus/sarocu.superplus.plist"
